# cms-xcache/cron.d/move_core
22 3,9,15,21 * * * root /usr/local/sbin/move_core > /dev/null 2>&1
