# cms-xcache/cron.d/refresh_proxy
39 3,9,15,21 * * * root /usr/local/sbin/refresh_proxy > /dev/null 2>&1
