# stash-origin/cron.d/stash-origin-auth.cron
*/15 * * * * /usr/libexec/xcache/authfile-update --origin
