# generate_gridmap.cron
30 */2 * * * root /usr/local/sbin/grid-mapfile.ligo-cvmfs.py > /dev/null 2>&1
