# xcache/cron.d/xcache-reporter.cron
*/5 * * * * * /usr/libexec/xcache/xcache-reporter --one-shot