# xcache/cron.d/xcache-reporter.cron
# OSG monitoring cron that respects the opt-out XRootD configuration
*/5 * * * * * root [ -n "$DISABLE_OSG_MONITORING" ] || /usr/libexec/xcache/xcache-reporter --one-shot
