# cms-xcache/systemd/system/cmsd@.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/cmsd -l /var/log/xrootd/cmsd.log -c /etc/xrootd/%i-cmsd.cfg -k fifo -s /var/run/xrootd/cmsd-%i.pid -n %i
