# stash-cache/cron.d/stash-cache-renew-proxy.cron
0 0 * * * xrootd /usr/libexec/xcache/renew-proxy
