# stash-cache/cron.d/stash-cache-authfile.cron
*/15 * * * * xrootd . /etc/xrootd-environment; /usr/libexec/xcache/authfile-update --cache
