# stash-cache/cron.d/stash-cache-auth.cron
*/4 * * * * /usr/libexec/xcache/authfile-update