# cms-xcache/supervisord.d/10-cms-xcache-cmsd.conf
[program:cms-xcache-cmsd]
command=/usr/bin/cmsd -l /var/log/xrootd/cmsd.log -c /etc/xrootd/xcache-cmsd.cfg -k fifo -s /var/run/xrootd/cmsd-xcache.pid -n xcache
user=xrootd
autorestart=true
