# cms-xcache/supervisord.d/10-cms-xcache.conf
[program:cms-xcache]
command=/usr/bin/xrootd -l /var/log/xrootd/xrootd.log -c /etc/xrootd/xcache.cfg -k fifo -s /var/run/xrootd/xrootd-xcache.pid -n xcache
user=xrootd
autorestart=true
